# output/goserver.service
[Unit]
Description=Golang Web Server

[Service]
ExecStart=/etc/init.d/goserver

[Install]
WantedBy=multi-user.target
